# legacyweb/conf/etc/cron.d/nablnews
34 12 * * * /home/ec2-user/bin/sendnews.sh > /home/ec2-user/news.log